# data/events/windows-security/4700.json
{
  "id": "4700",
  "source": "windows_security",
  "priority": "P3",
  "applicable_version": "Windows Server 2008 and later where the audit subcategory is enabled",
  "last_reviewed": "2026-07-22",
  "category": "ScheduledTask",
  "name": "A scheduled task was enabled",
  "definition": "Windows Security Event ID 4700 records that a scheduled task was enabled.",
  "trigger_scenarios": "Microsoft documents that this event is generated every time a scheduled task is enabled.",
  "key_fields": [
    {
      "field": "Task Name",
      "explanation": "The enabled task path. Microsoft recommends monitoring critical tasks that should never be enabled."
    },
    {
      "field": "Task Content",
      "explanation": "Task XML content that can expose command, arguments, trigger, and principal context."
    },
    {
      "field": "Subject Logon ID",
      "explanation": "Hexadecimal logon session for the account enabling the task."
    }
  ],
  "false_positives": [
    "Administrators may enable tasks during maintenance or application deployment.",
    "Vendor updaters can re-enable their own disabled tasks."
  ],
  "related_event_ids": [
    "4698",
    "4699",
    "4701",
    "4702"
  ],
  "attck_mapping": [
    {
      "tactic_id": "TA0003",
      "tactic_name": "Persistence",
      "technique_id": "T1053.005",
      "technique_name": "Scheduled Task/Job: Scheduled Task",
      "source_url": "https://attack.mitre.org/techniques/T1053/005/"
    }
  ],
  "detection_notes": "Use threshold=1 for critical Task Name values that should never be enabled, as Microsoft recommends. Task Name in the Task Scheduler Library root such as \\Updater or Task Content running powershell.exe from C:\\Users\\Public supports T1053.005 Scheduled Task when an actor reactivates persistence. Correlate Subject Logon ID with 4624/4688 and compare Task Content against the approved baseline.",
  "kql_snippet": "SecurityEvent\n| where EventID == 4700\n| project TimeGenerated, Computer, Account, TaskName, TaskContent",
  "spl_snippet": "index=wineventlog EventCode=4700 | table _time, host, Account_Name, Task_Name, Task_Content",
  "sample_log": "EventID: 4700\nSubject Logon ID: 0x42a91\nTask Name: \\Updater\nTask Content: <Command>C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe</Command>",
  "source_url": "https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4700",
  "technical_metadata": {
    "provider": "Microsoft-Windows-Security-Auditing",
    "channel": "Security",
    "level": "Information",
    "audit_keyword": "Audit Success"
  },
  "queries": [
    {
      "language": "kql",
      "title": "Retrieve Windows Security Event ID 4700",
      "query": "SecurityEvent\n| where EventID == 4700\n| project TimeGenerated, Computer, EventID, RenderedDescription\n| order by TimeGenerated desc",
      "data_source": "Microsoft Sentinel SecurityEvent table",
      "assumptions": [
        "The Windows Security Events connector populates SecurityEvent.",
        "Map raw XML fields to the event-specific fields described on this page."
      ],
      "source_url": "https://learn.microsoft.com/en-us/azure/azure-monitor/reference/tables/securityevent"
    },
    {
      "language": "spl",
      "title": "Retrieve Windows Security Event ID 4700",
      "query": "index=wineventlog source=\"WinEventLog:Security\" EventCode=4700\n| table _time, host, EventCode, Message",
      "data_source": "Splunk Windows Security event data",
      "assumptions": [
        "Index, source, and parsed field names vary by deployment.",
        "Apply the documented field conditions after validating local field extraction."
      ]
    },
    {
      "language": "xql",
      "title": "Retrieve Windows Security Event ID 4700",
      "query": "dataset = microsoft_windows_raw\n| filter edr_event_id = 4700",
      "data_source": "Cortex XDR microsoft_windows_raw",
      "assumptions": [
        "The Windows Security channel is ingested into Cortex XDR.",
        "Inspect raw provider fields because parsed names vary by collection path."
      ],
      "source_url": "https://docs-cortex.paloaltonetworks.com/r/Cortex-XDR/Cortex-XDR-5.x-Documentation/XDR-Collector-datasets"
    },
    {
      "language": "sigma",
      "title": "Windows Security Event ID 4700",
      "query": "title: Windows Security Event ID 4700\nstatus: experimental\nlogsource:\n    product: windows\n    service: security\ndetection:\n    selection:\n        EventID: 4700\n    condition: selection\nfalsepositives:\n    - Legitimate administrative or system activity\nlevel: informational",
      "data_source": "Sigma Windows Security log source",
      "assumptions": [
        "This rule retrieves the event and is not a standalone malicious-behavior detection.",
        "Add the event-specific conditions and correlations documented on this page."
      ],
      "source_url": "https://github.com/SigmaHQ/sigma-specification/blob/main/specification/sigma-rules-specification.md"
    }
  ],
  "value_references": [
    {
      "field": "EventID",
      "value": "4700",
      "meaning": "A scheduled task was enabled.",
      "security_relevance": "Interpret EventID with the other documented fields and correlated events; this value alone is not proof of compromise.",
      "source_url": "https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4700"
    },
    {
      "field": "TaskName",
      "value": "\\TaskName",
      "meaning": "The full scheduled-task path that was enabled.",
      "security_relevance": "Interpret TaskName with the other documented fields and correlated events; this value alone is not proof of compromise.",
      "source_url": "https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4700"
    }
  ],
  "faqs": [
    {
      "question": "What does Windows Event ID 4700 mean?",
      "answer": "A previously disabled scheduled task was enabled."
    },
    {
      "question": "Does Event ID 4700 always indicate malicious activity?",
      "answer": "No. Maintenance, deployment, and administrative workflows commonly enable tasks."
    },
    {
      "question": "How should Event ID 4700 be investigated?",
      "answer": "Inspect TaskName and TaskContent, then resolve SubjectLogonId to the initiating logon and process."
    }
  ],
  "sources": [
    {
      "title": "4700: A scheduled task was enabled",
      "url": "https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4700",
      "publisher": "Microsoft",
      "source_type": "vendor"
    },
    {
      "title": "Scheduled Task/Job: Scheduled Task (T1053.005)",
      "url": "https://attack.mitre.org/techniques/T1053/005/",
      "publisher": "MITRE ATT&CK",
      "source_type": "attack"
    },
    {
      "title": "XDR Collector datasets",
      "url": "https://docs-cortex.paloaltonetworks.com/r/Cortex-XDR/Cortex-XDR-5.x-Documentation/XDR-Collector-datasets",
      "publisher": "Palo Alto Networks",
      "source_type": "vendor"
    },
    {
      "title": "Sigma Rules Specification",
      "url": "https://github.com/SigmaHQ/sigma-specification/blob/main/specification/sigma-rules-specification.md",
      "publisher": "SigmaHQ",
      "source_type": "detection-rule"
    }
  ]
}
